# bot/bot.timer
[Unit]
Description=Daily run of bot at 10 AM ET

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
